# linux-cfg/var/spool/cron/crontab/banshee
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.P0vRpb/crontab installed on Mon Oct 14 13:36:55 2019)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# Edit this file to introduce tasks to be run by cron.
# 
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
# 
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').
# 
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
# 
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
# 
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command
MAILTO=""

*/30 * * * * /usr/local/ban/scripts/refreshFreeDns.sh
0 2 * * *    /usr/local/ban/scripts/updateListPackages.sh /usr/local/git/banserver/list-packages.txt > /usr/local/git/banserver/list-packages.txt
0 2 * * *    dpkg -l > /usr/local/git/banserver/list-dpkg.txt
#0 2 * * * /home/pi/scripts/gitBackupDrive.sh >> /var/log/banpi.log 2>&1

